# cookbooks/sample/templates/default/flaskapi.conf.erb
description "Flask Movie API"
author "Tim Palko"

start on runlevel [2345]

script 

	/bin/sleep <%= node[:service_delay_sec] %>
	cd <%= node[:base_folder] %>/services/api
	python flask_api.py

end script