# cookbooks/sample/templates/default/frontend.conf.erb
description "Movie Front End"
author "Tim Palko"

start on runlevel [2345]

script 
	
	/bin/sleep <%= node[:service_delay_sec] %>
	cd <%= node[:base_folder] %>/webapp
	python front_end.py

end script